# rc.local
#!/bin/bash
sleep 30
cd /home/oak/PhalaServerMonitor
./start.sh
exit 0