# rc.local
#!/bin/bash
sleep 30
cd ~/PhalaServerMonitor
./start.sh
exit 0